# src/rootfs.base/etc/systemd/system/container-getty@.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noreset --noclear --autologin $LOGINUSER - xterm-vt220
EnvironmentFile=/etc/default/kvdi
